# home/.config/systemd/user/stuffs.target
[Unit]
Description=My stuffs

[Install]
WantedBy=default.target
